# wcs_snapshot/debian/wcs.cron.d
MAILTO=root
LANG=C.UTF-8

* * * * * wcs /usr/bin/wcs-manage cron
